# nginx-gunicorn-flask/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn hello:app -w 20 -b localhost:5000
directory=/deploy/app
